# upcoming/bookworm/amd64/eth-node-testnet-besu-lighthouse/1.0.0-1/src/debian/postinst
#!/bin/bash

systemctl enable eth-node-testnet-besu.service
systemctl start eth-node-testnet-besu.service

systemctl enable eth-node-testnet-lighthouse-besu.service
systemctl start eth-node-testnet-lighthouse-besu.service