# upcoming/noble/all/eth-node-regtest/0.0.4-1/src/debian/clients/lighthouse-validator/eth-node-lighthouse-validator-regtest.service
[Unit]
Description=validator service file for eth-node-lighthouse for network: regtest
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/lib/eth-node-regtest/lighthouse-validator/run-lighthouse-validator-service.sh $DEBCRAFTER_EXTRA_SERVICE_ARGS
User=eth-node-lighthouse-val-regtest
Group=eth-node-lighthouse-val-regtest
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
RuntimeDirectory=eth-node-lighthouse-validator-regtest
RuntimeDirectoryMode=750
# no need to specify, these come from debcrafter
User=eth-node-lighthouse-val-regtest
# NoNewPrivileges=true
# ProtectHome=true
# PrivateTmp=true
# PrivateDevices=true

# additional flags not specified by debcrafter
CapabilityBoundingSet=
IPAddressDeny=none
LockPersonality=true
PrivateIPC=true
PrivateUsers=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ReadWritePaths=/var/lib/eth-node-regtest/lighthouse-validator
ReadOnlyPaths=/var/lib/eth-node-regtest
RemoveIPC=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK 

RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
WorkingDirectory=/var/lib/eth-node-regtest/lighthouse-validator


[Install]
WantedBy=multi-user.target
